# config/systemd/ron-deadline-check.timer
[Unit]
Description=Ron Clanker - Check FPL deadlines every 2 hours

[Timer]
OnCalendar=*-*-* 00/2:00:00
Persistent=true
RandomizedDelaySec=120

[Install]
WantedBy=timers.target
