# config/systemd/ron-data-collection.timer
[Unit]
Description=Ron Clanker - Daily FPL data collection at 02:30

[Timer]
OnCalendar=*-*-* 02:30:00
Persistent=true

[Install]
WantedBy=timers.target
