# config/systemd/ron-daily-scout.timer
[Unit]
Description=Ron Clanker - Daily intelligence at 03:00

[Timer]
OnCalendar=*-*-* 03:00:00
Persistent=true

[Install]
WantedBy=timers.target
